# services/plantuml.socket
[Unit]
Description=Socket activation for PlantUML-Server

[Socket]
ListenStream=17080
Accept=no

[Install]
WantedBy=default.target
